# services/rfx-services.target
#
#  This file is part of systemd rfx-services
#

[Unit]
Description=rfx-services
After=network-online.target docker.service
