# services/rfx-services.service.in
[Unit]
Description=Start rfx-services
After=docker.service docker.socket

[Service]
Type=oneshot
ExecStart=@prefix@/etc/init.d/rfx-services.sh start

[Install]
WantedBy=multi-user.target